# deploy/launchd.example.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS schedule: Mon/Wed/Sat at 09:00.
  Install:
    cp deploy/launchd.example.plist ~/Library/LaunchAgents/com.parkrun-monitoring.plist
    (edit WorkingDirectory below, then)
    launchctl load ~/Library/LaunchAgents/com.parkrun-monitoring.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.parkrun-monitoring</string>
    <key>WorkingDirectory</key>
    <string>/REPLACE/WITH/PATH/TO/parkrun-monitoring</string>
    <key>ProgramArguments</key>
    <array>
        <string>.venv/bin/parkrun-monitoring</string>
        <string>sync</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>6</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    </array>
    <key>StandardOutPath</key>
    <string>data/sync.log</string>
    <key>StandardErrorPath</key>
    <string>data/sync.log</string>
</dict>
</plist>
